# systemd/claude-daemon.service
[Unit]
Description=Claude Discord Plugin Daemon
After=network.target

[Service]
Type=forking
User=shankyram
Environment="HOME=/home/shankyram"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStart=/usr/bin/tmux new-session -d -s claude-daemon 'bash /data/tools/trading-claude/scripts/claude_daemon.sh'
ExecStop=/usr/bin/tmux kill-session -t claude-daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target